# packaging/linux/preremove.sh
#!/bin/bash
systemctl stop comenqd.service
systemctl disable comenqd.service
